# deploy/systemd/ocpctl-worker.service
[Unit]
Description=OCPCTL Worker
After=network.target postgresql.service
Wants=network-online.target

[Service]
Type=simple
User=ocpctl
WorkingDirectory=/opt/ocpctl
EnvironmentFile=/etc/ocpctl/worker.env

# Ensure installer binaries are available before starting (runs as root with + prefix)
ExecStartPre=+/opt/ocpctl/scripts/ensure-installers.sh

ExecStart=/opt/ocpctl/current/ocpctl-worker
Restart=on-failure
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=ocpctl-worker

# Graceful shutdown configuration
# Give jobs up to 1 hour to complete on shutdown
TimeoutStopSec=3600
# Send SIGTERM to main process first, only SIGKILL children if timeout exceeded
KillMode=mixed
# Use SIGTERM for graceful shutdown (default, but explicit for clarity)
KillSignal=SIGTERM

# Resource limits
LimitNOFILE=65536
LimitNPROC=4096

# Security
NoNewPrivileges=true
PrivateTmp=true

[Install]
WantedBy=multi-user.target
